# firefox-jenkins-slave/xvfb.sv.conf
[program:xvfb]
command=/bin/bash -c "Xvfb :10 -ac"