# jenkins-slave/jenkins-slave.sv.conf
[program:jenkins-slave]
command=/bin/bash -c "runslave"
