# elasticsearch/supervisord.conf
[program:elasticsearch]
command=/bin/bash -c "cd /opt/elasticsearch-0.90.5 && bin/elasticsearch -f"
redirect_stderr=true
