# elasticsearch/supervisord.conf
[supervisord]
nodaemon=true

[program:elasticsearch]
command=/opt/elasticsearch-0.90.5/bin/elasticsearch
