# dood/docker-daemon.sv.conf
[program:docker]
command=/bin/bash -c "docker -d"
