# jenkins/jenkins.sv.conf
[program:jenkins]
command=java -jar /var/lib/jenkins/jenkins.war
redirect_stderr=true
