# elasticsearch-1.x/elasticsearch.sv.conf
[program:elasticsearch]
command=/bin/bash -c "cd /opt/elasticsearch-1.0.1 && bin/elasticsearch -f"
redirect_stderr=true
