# dind-jenkins-slave/dnsmasq.sv.conf
[program:dnsmasq]
command=/bin/bash -c "/usr/sbin/dnsmasq -k"
