# dind-jenkins-slave/docker.sv.conf
[program:docker]
command=/bin/bash -c "wrapdocker && docker -d"
