# dind-jenkins-slave/supervisord.conf
[supervisord]
nodaemon=true

[program:docker]
command=/bin/bash -c "wrapdocker && docker -d"

[program:jenkins-slave]
command=/bin/bash -c "runslave"