# mongo/mongo.sv.conf
[program:mongo]
command=/opt/mongodb/bin/mongod
redirect_stderr=true
