# mongo/supervisord.conf
[supervisord]
nodaemon=true

[program:mongo]
command=/opt/mongodb/bin/mongod
